# rootfiles/quadlets/funkwhale/units/funkwhale-redis.container
[Container]
Image=redis:7-alpine
Pod=funkwhale.pod
EnvironmentFile=/etc/funkwhale.env
EnvironmentFile=%h/.config/funkwhale-secrets.env
Volume=funkwhale-redis-data:/data

[Service]
Restart=always
