# modules/funkwhale/units/funkwhale-api.container
[Unit]
Requires=funkwhale-postgres.service funkwhale-redis.service
After=funkwhale-postgres.service funkwhale-redis.service

[Container]
Image=funkwhale/api:2.0.3
Pod=funkwhale.pod
EnvironmentFile=/etc/funkwhale/funkwhale.env
EnvironmentFile=%h/.config/funkwhale-extra.env
Exec=gunicorn
Volume=%h/.local/share/funkwhale/volumes/funkwhale-music:/music:z,ro
Volume=%h/.local/share/funkwhale/volumes/funkwhale-media:/media:z
Volume=funkwhale-static:/static

[Service]
Restart=always
